# crontab/update-dns.crontab
*/5 * * * * BASH_ENV=/home/pi/.profile /home/pi/update-dns.sh >> /home/pi/update-dns.log 2>&1